# psutild.service
[Service]
ExecStart=/home/student/PycharmProjects/psutilDaemonTask/task.py
